# docker-s6/root/etc/services.d/trafegodns/finish
#!/usr/bin/with-contenv bash

echo
echo "-------------------------------------"
echo "TráfegoDNS has exited with status $?. Restarting..."
echo "-------------------------------------"
echo
exit 0