# docker-s6/root/etc/services.d/trafegodns/finish
#!/usr/bin/with-contenv sh

echo "TráfegoDNS has exited with status $?. Restarting..."
exit 0