# docker-s6/root/etc/services.d/trafegodns/run
#!/usr/bin/with-contenv sh

echo "Starting TráfegoDNS..."
cd /app
exec s6-setuidgid abc node src/app.js