# docker-s6/root/etc/services.d/webui/run
#!/command/with-contenv bash
# /etc/services.d/webui/run

echo
echo "-------------------------------------"
echo "Starting TráfegoDNS Web UI..."
echo "-------------------------------------"
echo

# Check if Web UI is enabled
if [ "$ENABLE_WEB_UI" != "true" ]; then
  echo "Web UI is disabled, set ENABLE_WEB_UI=true to enable"
  echo "Exiting Web UI service"
  exit 0
fi

# Set environment variables for web ui
export ENABLE_WEB_UI=true
export WEB_UI_PORT=${WEB_UI_PORT:-8080}

# Check for authentication
if [ -n "$WEB_UI_USERNAME" ] && [ -n "$WEB_UI_PASSWORD" ]; then
  echo "Web UI authentication enabled"
else
  echo "Warning: Web UI authentication not configured"
  echo "Consider setting WEB_UI_USERNAME and WEB_UI_PASSWORD environment variables"
fi

# Output Web UI information
echo "Web UI will be available at http://localhost:${WEB_UI_PORT}"
echo "Web UI files are served from: /app/webui/build"
echo "-------------------------------------"

# The Web UI is integrated with the main application
# It will be served by the Express server in app.js when ENABLE_WEB_UI=true
# This service just provides logging and keeps the service running

# Keep the service running
while true; do
  sleep 3600
done